# .docker/development/files/etc/cron.d/www-data

# Consume messenger queue every minute
* * * * *    /usr/local/bin/php bin/console messenger:consume --time-limit 59 async &>/dev/null
